# multipath/mpcount.service
[Unit]
Description=Multipath Count Service

[Service]
Type=forking
ExecStart=/etc/udev/scripts/udevmpath.sh
