# systemd/mpathcount.socket
[Unit]
Description=Multipath scanner kick socket

[Socket]
ListenFIFO=/var/run/mpathcount.sock
PipeSize=4096

[Install]
WantedBy=multi-user.target
