# systemd/usb-scan.socket
[Unit]
Before=systemd-udevd

[Socket]
ListenFIFO=/var/run/usb-scan.sock
PipeSize=4096

[Install]
WantedBy=multi-user.target
